# deploy/systemd/pouchcare-deploy.timer
[Unit]
Description=PouchCare auto-deploy poller (every ~60s)

[Timer]
OnBootSec=45
OnUnitActiveSec=60
AccuracySec=15
Unit=pouchcare-deploy.service

[Install]
WantedBy=timers.target
